# supervisord.conf
[supervisord]
nodaemon=true

[program:phpfpm]
command=/bin/bash -c "/usr/sbin/php-fpm --fpm-config /etc/php/fpm/php-fpm.conf -F"
